# target/supervisor_restic.ini
[supervisord]
nodaemon=true

[program:cron]
startsecs=0
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
command=/go/bin/start_cron.sh

[program:restic_backup]
startsecs=0
autostart=false
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=/go/bin/restic-runner --repo %(ENV_RR_REPO)s --set %(ENV_RR_SET)s backup

[program:restic_cleanup]
startsecs=0
autostart=false
autorestart=unexpected
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=/go/bin/restic-runner --repo %(ENV_RR_REPO)s expire

